# bin/harvest.cron
0 23 * * * * /path/to/harvest http://www.ecb.europa.eu/stats/eurofxref/eurofxref-hist-90d.xml <database_string>
